# crontab/crontab

0 * * * * /usr/local/bin/python3 /application/session_cleanup.py
