# crontab/crontab

0 * * * * /usr/local/bin/python3 /application/update_database.py
